# rpled.service
[Unit]
Description=Controll LED

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/bin/python /home/pi/rpled/rpled.py
User=pi
Group=pi
